# files/timers/cache-updates.timer
[Timer]
OnBootSec=1h
OnUnitInactiveSec=1h
Persistent=true
RandomizedDelaySec=30m
