# files/timers/cache-updates.timer
[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=30m
