# files/timers/clean-cache.timer
[Timer]
OnBootSec=1m
OnUnitInactiveSec=24h
Persistent=true
RandomizedDelaySec=30m
