# files/timers/clean-cache.timer
[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=3h
