# files/timers/clean-pacaur-cache.service
[Service]
Type=oneshot
ExecStart=/bin/rm -rf /home/makepkg/.cache/pacaur
