# files/timers/localepurge.service
[Service]
Type=oneshot
ExecStart=/usr/bin/localepurge
