# files/timers/pacman-fy.timer
[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=30m
