# files/timers/reflector.timer
[Timer]
OnCalendar=hourly
RandomizedDelaySec=30m
