# files/timers/reflector.service
[Service]
Type=oneshot
ExecStart=/bin/reflector --sort delay --save /etc/pacman.d/mirrorlist
