# files/timers/reflector.timer
[Timer]
OnCalendar=daily
RandomizedDelaySec=6h
