# files/powerdown.service
[Unit]
Description=Laptop battery savings

[Service]
Type=oneshot
ExecStart=/usr/local/bin/powerdown

[Install]
WantedBy=battery.target
