# files/powerup.service
[Unit]
Description=Disable laptop battery savings

[Service]
Type=oneshot
ExecStart=/usr/local/bin/powerup

[Install]
WantedBy=ac.target
